# config/com.oasis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this launch agent -->
    <key>Label</key>
    <string>com.oasis.organizer</string>

    <!-- Launch the OASIS app (which has Full Disk Access) -->
    <!-- Note: __HOME__ is replaced by install.sh with actual home directory -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>-a</string>
        <string>__HOME__/.local/share/oasis/OASIS.app</string>
    </array>

    <!-- Run daily at midnight -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Also run when the agent is first loaded (on login or after install) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Standard output and error logging -->
    <key>StandardOutPath</key>
    <string>__HOME__/.config/oasis/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.config/oasis/stderr.log</string>
</dict>
</plist>
